# api/docker/supervisor/event.ini
[program:event]
directory=/srv/app
command=/bin/sh -c '/srv/app/bin/console  messenger:consume async -vvv --limit=100 --time-limit=3600 --no-reset'
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
user=1000
